# ops/plutus-homelab-user.service
[Unit]
Description=Plutus homelab storefront API (:8030)
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/ai-workspace/plutus
EnvironmentFile=%h/ai-workspace/plutus/.env.homelab
ExecStart=%h/ai-workspace/plutus/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8030
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target